# scripts/daily-update.crontab
* * * * * root /home/biobuddy/scripts/daily-update.sh >> /var/log/cron.log 2>&1
