# app/daily-update/daily-update.crontab
0 23 * * * root /home/biobuddy/app/daily-update/daily-update.sh >> /var/log/cron.log 2>&1
